# deploy/junshi-api.service
# 军师 · 后端 systemd 服务单元
# 安装：拷到 /etc/systemd/system/junshi-api.service → 改路径/用户 →
#       systemctl daemon-reload && systemctl enable --now junshi-api
# 日志：journalctl -u junshi-api -f
[Unit]
Description=Junshi API (Fastify + Prisma)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=junshi
WorkingDirectory=/opt/junshi/server
Environment=TZ=Asia/Shanghai
# 环境变量从 WorkingDirectory 下的 .env 读取（应用内置 dotenv 自动加载）。
# 如偏好 systemd 注入，可改用： EnvironmentFile=/opt/junshi/server/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
# 适度加固
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
